# root/etc/services.d/mysql/run
#!/usr/bin/with-contenv bash
/usr/bin/mysqld_safe --datadir=/data --defaults-file=/config/my.cnf